# systemd/iron-agent.service
[Unit]
Description=IRON Secure Reverse Tunnel Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/iron/iron.py agent --config /etc/iron/agent.json
Restart=always
RestartSec=3
LimitNOFILE=1048576
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/etc/iron

[Install]
WantedBy=multi-user.target
